# packages/marcopolo-deployer/marcopolo-deployer-0.1.6.tar.gz/marcopolo-deployer-0.1.6/daemon/systemv/marcodeployerd
#!/bin/sh

### BEGIN INIT INFO
# Provides:          marcodeployerd
# Required-Start:    $remote_fs $syslog marcod polod
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: The marcopolo deployer
# Description:       Marcodeployer is a deploying tool based on
#                    MarcoPolo
### END INIT INFO

# Author: Diego Martin <martinarroyo@usal.es>
#

#FROM: https://twistedmatrix.com/trac/attachment/ticket/3434/twistdplugin

# Do NOT "set -e"

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
NAME="marcodeployerd"
DESC="The marcopolo discoverer"
DAEMON=$(which marcodeployerd)

PIDFILE=/var/run/${NAME}.pid
LOGFILE=/var/log/${NAME}.log

SERVICE_NAME=${NAME}
DAEMON_OPTS=

if [ ! -x $DAEMON ]; then
  echo "ERROR: Can't execute $DAEMON."
  exit 1
fi

# Read configuration variable file if it is present
[ -r /etc/default/$NAME ] && . /etc/default/$NAME

# Load the VERBOSE setting and other rcS variables
. /lib/init/vars.sh

# Define LSB log_* functions.
# Depend on lsb-base (>= 3.2-14) to ensure that this file is present
# and status_of_proc is working.
. /lib/lsb/init-functions

start_service() {
  echo -n " * Starting $NAME... "
  start-stop-daemon -Sq --make-pidfile -p $PIDFILE -x $DAEMON -- $DAEMON_OPTS &>/dev/null
  e=$?
  if [ $e -eq 1 ]; then
    echo "already running"
    return
  fi

  if [ $e -eq 255 ]; then
    echo "couldn't start :("
    return
  fi

  echo "done"
}

stop_service() {
  echo -n " * Stopping $NAME... "
  start-stop-daemon -Kq -R 10 -p $PIDFILE
  e=$?
  if [ $e -eq 1 ]; then
    echo "not running"
    return
  fi

  echo "done"
}

case "$1" in
  start)
    start_service
    ;;
  stop)
    stop_service
    ;;
  restart)
    stop_service
    start_service
    ;;
  status)
    status_of_proc "$DAEMON" "$NAME" && exit 0 || exit $?
    ;;
  *)
    echo "Usage: /etc/init.d/$SERVICE_NAME {start|stop|restart|status}" >&2
    exit 1   
    ;;
esac

exit 0