# packages/marcopolo-deployer/marcopolo-deployer-0.1.6.tar.gz/marcopolo-deployer-0.1.6/daemon/systemd/marcodeployerd.service
[Unit]
Description=A deployment utility
Requires=polod.service
After=polod.service

[Service]

ExecStart=/usr/bin/marcodeployerd &> /dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/marcodeployerd.pid

[Install]
WantedBy=multi-user.target
