# packages/marcobootstrap.slave/marcobootstrap.slave-0.0.18.tar.gz/marcobootstrap.slave-0.0.18/daemon/systemd/mbootstrapslad.service
[Unit]
Description=The marcobootstrap client
After=marco.service
After=polo.service

[Service]

ExecStart=/usr/bin/mbootstrapslad &>/dev/null

Type=simple
User=root
Group=root
PIDFile=/var/run/mbootstrapslad.pid

[Install]
WantedBy=multi-user.target
