# packages/peekaboo/peekaboo-0.8.5.tar.gz/peekaboo-0.8.5/contrib/iseeyou.service
[Unit]
Description=iseeyou
After=network.target
Before=peekaboo.service

[Service]
Restart=on-failure
RestartSec=10
Type=simple
User=iseeyou
Group=iseeyou
ExecStart=/usr/bin/iseeyou

[Install]
WantedBy=multi-user.target
